# ops/launchd/com.norikane.hp41238.staging-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.norikane.hp41238.staging-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/norikene_satoshi/projects/norikane_satoshi_HP/.codex-worktrees/staging-live-41238/scripts/hp-41238-sync-staging.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/norikene_satoshi/.local/share/hp-41238/logs/staging-sync.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/norikene_satoshi/.local/share/hp-41238/logs/staging-sync.launchd.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/norikene_satoshi</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/Users/norikene_satoshi/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
